# vm-support/systemd/timer-daily.timer
[Unit]
Description=Daily Timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=timer-daily.target

[Install]
WantedBy=multi-user.target
